# e-ink_code/server/fun_aggregator/deploy/fun-aggregator.service
# Install on Raspberry Pi (repo checkout at /opt/e-ink-display):
#   sudo apt install python3-venv
#   cd /opt/e-ink-display/e-ink_code/server/fun_aggregator
#   python3 -m venv .venv && .venv/bin/pip install -r requirements.txt
#   .venv/bin/python run_server.py   # listens on 0.0.0.0:8081
#   sudo cp deploy/fun-aggregator.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now fun-aggregator
#
# Alternate layout: rsync server/fun_aggregator/ to /opt/fun-aggregator/ and edit
# WorkingDirectory + ExecStart below to /opt/fun-aggregator and its .venv.
#
# Optional env file (FUN_API_KEY strongly recommended for any network exposure):
#
#   FUN_API_KEY=<long random secret>   # mandatory if FUN_REQUIRE_API_KEY is on
#   FUN_REQUIRE_API_KEY=1              # fail fast at startup if FUN_API_KEY unset (recommended on Pi)
#
#   FUN_RATE_LIMIT_SCREEN=60/minute    # per-client-IP limit on GET /v1/fun/screen
#   FUN_RATE_LIMIT_BATCH=40/minute    # per-client-IP limit on GET /v1/fun/facts/batch
#
#   FUN_DEVICE_STORE=data/devices.json # device roster (UUID ↔ friendly_name); parent dir created if missing
#   FUN_SPECIAL_STORE=data/special_messages.json  # FIFO special slides per device + named groups JSON
#   FUN_ADMIN_API_KEY=<secret>  # POST /v1/admin/special (header X-Fun-Admin-Key); omit to disable admin
#
#   REFRESH_SECONDS=900
#
#   Optional upstream fact APIs (see server/README.md). Built-in defaults: meowfacts + uselessfacts.jsph
#   when FUN_FACT_SOURCES_JSON is unset and per-URL env vars are empty.
#   Live check: cd /opt/e-ink-display/e-ink_code/server/fun_aggregator && .venv/bin/python smoke_fact_upstreams.py
#
# Example:
#   sudo nano /etc/fun-aggregator.env

[Unit]
Description=Fun facts aggregator for e-ink displays
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/e-ink-display/e-ink_code/server/fun_aggregator
EnvironmentFile=-/etc/fun-aggregator.env
# Fun API listens on port 8081 (see run_server.py; override with FUN_LISTEN_PORT).
Environment=FUN_LISTEN_HOST=0.0.0.0
Environment=FUN_LISTEN_PORT=8081
ExecStart=/opt/e-ink-display/e-ink_code/server/fun_aggregator/.venv/bin/python run_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
